# Substr8HUD/systemd/substr8-hud-console@.service.in
[Unit]
Description=SUBSTR8 HUD shared physical console for %I
After=systemd-user-sessions.service

[Service]
Type=simple
Environment=TERM=linux
EnvironmentFile=-/etc/substr8-hud/console-%i.conf
ExecStart=@bindir@/substr8-hud-console system-attach
Restart=always
RestartSec=1
KillMode=process
StandardInput=tty-force
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty%i
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no

[Install]
WantedBy=multi-user.target
